# .config/systemd/user/test.service
[Unit]
Description=Update the plocate database

[Service]
Type=oneshot
ExecStart=%h/bin/test.sh
LimitNOFILE=131072
IOSchedulingClass=idle
