# .config/Configs/ServiceFiles/dmlock.service
[Unit]
Description=DM Lock before sleep
Before=sleep.target

[Service]
ExecStart=/etc/systemd/system/lockscreen.sh

[Install]
WantedBy=sleep.target
